# packages/usbnb/New folder/New folder/packages/deb/prerm
#!/bin/bash
set -e

# Debian prerm script
# Called before package removal

SERVICE_NAME="networkbuster"

# Stop the service if running
if systemctl is-active --quiet ${SERVICE_NAME} 2>/dev/null; then
    systemctl stop ${SERVICE_NAME} || true
fi

# Disable the service
if systemctl is-enabled --quiet ${SERVICE_NAME} 2>/dev/null; then
    systemctl disable ${SERVICE_NAME} || true
fi

exit 0
